# auto-config/udev/99-hilink.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="12d1", RUN+="/usr/local/bin/hilink-autoconf.sh"